# backend/opendeploy.service
[Unit]
Description=OpenDeploy Hosting Dashboard
After=network-online.target avahi-daemon.service
Wants=network-online.target avahi-daemon.service

[Service]
Type=simple
User=opendeploy
Group=opendeploy

ExecStart=/usr/local/bin/opendeploy --config /etc/opendeploy/config.yaml
Restart=always
RestartSec=5
StartLimitIntervalSec=200
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=opendeploy
Environment=opendeploy_ENV=production

NoNewPrivileges=yes

PrivateTmp=yes

ProtectSystem=strict

ReadWritePaths=/var/lib/opendeploy /var/log/opendeploy /etc/cloudflared /etc/nginx/sites-available /etc/nginx/sites-enabled

ProtectKernelLogs=yes
ProtectKernelTunables=yes

PrivateDevices=yes
MemoryDenyWriteExecute=yes

[Install]
WantedBy=multi-user.target